# pkgs/hostapd-edison-cfg/DEBIAN/postinst
#!/bin/sh

update-rc.d hostapd disable
systemctl disable udhcpd.service
